# nginx-sync/nginx-sync.timer
[Unit]
Description=nginx-sync timer

[Timer]
OnUnitActiveSec=5min
OnBootSec=5min

[Install]
WantedBy=timers.target
